# docker/supervisord.conf
[supervisord]
nodaemon=true

[group:plan]
programs=Active-plan,Upload-store-list,Store-download,Protocol-Store-download,Protocol-Batch-upload

[group:image]
programs=Dequeue-cos,Dequeue-calculation-task,Dequeue-calculation-result,Send-result-to-zft,Dequeue-replan-task

[group:replan]
programs=Replan-create,Push-engine

#检查结果下载
[group:report]
programs=Scene-download,Question-download

[group:plan-batch]
programs=plan-batch-excel-import,plan-batch-excel-import-fail-download

[program:php-fpm]
command=/usr/sbin/php-fpm7.2 --nodaemonize
autostart=true
autorestart=true
stopwaitsecs=60
stdout_logfile=/var/log/supervisor/php-fpm.log

[program:nginx]
command = /usr/sbin/nginx -g 'daemon off;'
startsecs=0
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stopasgroup=true
killasgroup=true

[program:env]
command=/application/env.sh
nodaemon=true
startretries=1
autostart=true
autorestart=false
startsecs=0

# 检查计划售点上传
[program:Upload-store-list]
command=/application/yii plan/upload-store-list
autostart=true
stopwaitsecs=60
process_name=plan_upload_store_list_%(process_num)s
numprocs=5

# 检查计划生成售点检查数据场景和问卷
[program:Active-plan]
command=/application/yii plan/active-plan
autostart=true
stopwaitsecs=60
process_name=plan_Active_plan_%(process_num)s
numprocs=5

# 检查计划已配置售点数据下载
[program:Store-download]
command=/application/yii plan/store-download
autostart=true
stopwaitsecs=60
process_name=plan_store_download_%(process_num)s
numprocs=5

# 检查计划已配置售点数据下载
[program:Protocol-Store-download]
command=/application/yii plan/protocol-store-download
autostart=true
stopwaitsecs=60
process_name=protocol_plan_store_download_%(process_num)s
numprocs=5

# 检查计划导入售点 成功数据下载
[program:Import-Success-download]
command=/application/yii plan/excel-store-download
autostart=true
stopwaitsecs=60
process_name=excel_store_download_%(process_num)s
numprocs=5

# 检查计划导入售点 失败数据下载
[program:Import-Fail-download]
command=/application/yii plan/excel-store-fail-download
autostart=true
stopwaitsecs=60
process_name=excel_store_fail_download_%(process_num)s
numprocs=5

# 检查计划导入售点 失败数据下载
[program:Plan-Store-Relation-Update]
command=/application/yii plan/plan-store-relation-update
autostart=true
stopwaitsecs=60
process_name=plan_store_relation_update_%(process_num)s
numprocs=5

# 批量创建检查计划
[program:Protocol-Batch-upload]
command=/application/yii plan/plan-batch-upload
autostart=true
stopwaitsecs=60
process_name=plan_batch_upload_%(process_num)s
numprocs=5

#将图像上传cos后，送图像识别
[program:Dequeue-cos]
command=/application/yii image-service/dequeue-cos
autostart=true
stopwaitsecs=60
process_name=image_service_dequeue_cos_%(process_num)s
numprocs=5

#走访完成创建任务，判断条件送引擎计算
[program:Dequeue-calculation-task]
command=/application/yii image-service/dequeue-calculation-task
autostart=true
stopwaitsecs=60
process_name=image_service_dequeue_calculation_task_%(process_num)s
numprocs=5

#接收引擎计算结果
[program:Dequeue-calculation-result]
command=/application/yii image-service/dequeue-calculation-result
autostart=true
stopwaitsecs=60
process_name=image_dequeue_calculation_result_%(process_num)s
numprocs=5

#拼凑条件字段将结果送给zft的脚本
[program:Send-result-to-zft]
command=/application/yii image-service/send-result-to-zft
autostart=true
stopwaitsecs=60
process_name=image_send_result_to_zft_%(process_num)s
numprocs=5

#重跑数据入规则引擎计算的脚本
[program:Dequeue-replan-task]
command=/application/yii image-service/dequeue-replan-task
autostart=true
stopwaitsecs=60
process_name=image_dequeue_replan_task_%(process_num)s
numprocs=5

# 检查结果引擎结果下载
[program:Scene-download]
command=/application/yii report/scene-download
autostart=true
stopwaitsecs=60
process_name=report_scene_download_%(process_num)s
numprocs=5

# 检查结果问卷结果下载
[program:Question-download]
command=/application/yii report/question-download
autostart=true
stopwaitsecs=60
process_name=report_question_download_%(process_num)s
numprocs=5

[program:Replan-create]
command=/application/yii replan/replan-create
autostart=true
stopwaitsecs=60
process_name=replan_create_%(process_num)s
numprocs=5

[program:Push-engine]
command=/application/yii replan/push-engine
autostart=true
stopwaitsecs=60
process_name=replan_push_engine_%(process_num)s
numprocs=2

# 长期协议门店关系批量导入
[program:plan-batch-excel-import]
command=/application/yii plan-batch/excel-import
autostart=true
stopwaitsecs=60
process_name=plan-batch-excel-import-%(process_num)s
numprocs=5

# 长期协议门店关系批量导入失败文件下载
[program:plan-batch-excel-import-fail-download]
command=/application/yii plan-batch/excel-import-fail-download
autostart=true
stopwaitsecs=60
process_name=plan-batch-excel-importfail-download-%(process_num)s
numprocs=1

# 所有下载统一出队
[program:download]
command=/application/yii base/dequeue
autostart=true
stopwaitsecs=60
process_name=download_%(process_num)s
numprocs=20

[program:Qc-manual-check-result-list-download]
command=/application/yii qc/manual-check-result-list-download
autostart=true
stopwaitsecs=60
process_name=manual_check_result_list_download_%(process_num)s
numprocs=1

# 问卷QC批量复核售点清单导入
[program:question-qc-survey-ignore-upload]
command=/application/yii question-qc/survey-ignore-upload
autostart=true
stopwaitsecs=60
process_name=question-qc-survey-ignore-upload-%(process_num)s
numprocs=1

# 问卷QC批量复核售点清单导入失败文件下载
[program:question-qc-fail-download]
command=/application/yii question-qc/fail-download
autostart=true
stopwaitsecs=60
process_name=question-qc-fail-download-%(process_num)s
numprocs=1